# startWPASupp.service
[Unit]
Description=startWPASupp

[Service]
ExecStart=/home/bin/startWPASupp

[Install]
WantedBy=multi-user.target

